# deploy/docmolder-db-backup.timer
[Unit]
Description=Daily DocMolder SQLite backup

[Timer]
OnCalendar=*-*-* 03:25:00
Persistent=true
Unit=docmolder-db-backup.service

[Install]
WantedBy=timers.target
