# deploy/docmolder-alertcheck.service
[Unit]
Description=DocMolder operational alert check
After=network-online.target docmolder.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/docmolder/app/deploy/alert-check.sh
