# deploy/docmolder-reconcile.timer
[Unit]
Description=Periodic DocMolder runtime reconciliation

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Persistent=true
Unit=docmolder-reconcile.service

[Install]
WantedBy=timers.target
