# deploy/docmolder-db-backup.service
[Unit]
Description=DocMolder SQLite backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/docmolder/app/deploy/backup-db.sh
