# deploy/docmolder-alertcheck.timer
[Unit]
Description=Periodic DocMolder operational alert check

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true
Unit=docmolder-alertcheck.service

[Install]
WantedBy=timers.target
